# memwatch.service
[Unit]
Description=Dynamic user memory watchdog
After=network.target

[Service]
Type=simple
Nice=-1
ExecStart=/usr/sbin/memwatch

[Install]
WantedBy=default.target
